# binaries/android/services/mongod/run
#!/data/data/com.termux/files/usr/bin/sh
exec mongod --dbpath "$HOME/mongodb-data" --bind_ip 127.0.0.1 2>&1
